# spotify_display.service
[Unit]
Description=Spotify Matrix Display Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mever
WorkingDirectory=/home/mever/spotify-matrix
Environment=PYTHONPATH=/home/mever/spotify-matrix:/home/mever/spotify-matrix/rpi-rgb-led-matrix/bindings/python:/home/mever/spotify-matrix/venv/lib/python3.11/site-packages:/usr/local/lib/python3.11/dist-packages
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/mever
ExecStart=/usr/bin/sudo -E /home/mever/spotify-matrix/venv/bin/python spotify_display_main.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target 