# spotify_display.service
[Unit]
Description=Spotify Matrix Display Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/mever/spotify-matrix
Environment=PYTHONPATH=/home/mever/spotify-matrix:/home/mever/spotify-matrix/rpi-rgb-led-matrix/bindings/python:/home/mever/spotify-matrix/venv/lib/python3.11/site-packages
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/mever/spotify-matrix/venv/bin/python spotify_display_main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target 